# agent/agent.service
[Unit]
Description=agent

[Service]
ExecStart=/usr/local/bin/agent.py